# devops/supervisord.conf
[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[supervisord]
nodaemon=true
loglevel=info
logfile=/var/log/supervisord.log
pidfile=/var/run/supervisord.pid

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[program:octane]
command=php /var/www/artisan octane:start --server=roadrunner --host=127.0.0.1 --port=8000 --workers=4 --max-requests=500
directory=/var/www
autostart=true
autorestart=true
startretries=3
user=www-data
stdout_logfile=/var/log/octane.log
stderr_logfile=/var/log/octane_err.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB
stdout_logfile_backups=3
stderr_logfile_backups=3
stopwaitsecs=10

[program:nginx]
command=nginx -g "daemon off;"
autostart=true
autorestart=true
startretries=3
stdout_logfile=/var/log/nginx_access.log
stderr_logfile=/var/log/nginx_error.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB
stdout_logfile_backups=3
stderr_logfile_backups=3

[program:laravel-worker]
command=php -d memory_limit=1G /var/www/artisan queue:work database --sleep=3 --tries=3 --max-time=3600 --memory=1024 --timeout=3600
directory=/var/www
autostart=true
autorestart=true
startretries=3
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=/var/log/laravel-worker.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
stopwaitsecs=60

[program:laravel-scheduler]
command=bash -c "while true; do php /var/www/artisan schedule:run && sleep 60; done"
directory=/var/www
autostart=true
autorestart=true
startretries=3
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=/var/log/laravel-scheduler.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3

[program:reverb]
command=php /var/www/artisan reverb:start --host=0.0.0.0 --port=8081
directory=/var/www
autostart=true
autorestart=true
startretries=3
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=/var/log/reverb.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
stopwaitsecs=10

# Ensure the log files are writable by www-data
[program:setup-permissions]
command=bash -c "touch /var/log/laravel-worker.log /var/log/laravel-scheduler.log /var/log/reverb.log /var/log/supervisord.log && chmod 664 /var/log/laravel-worker.log /var/log/laravel-scheduler.log /var/log/reverb.log /var/log/supervisord.log && chown www-data:www-data /var/log/laravel-worker.log /var/log/laravel-scheduler.log /var/log/reverb.log"
autostart=true
autorestart=false
startsecs=0
priority=1
